# overlay/etc/local.d/load_sysctl.start
#!/bin/bash

sysctl -p /etc/sysctl.d/*.conf
